# deploy/auto-translate-update.timer
[Unit]
Description=Auto-translate bot update timer

[Timer]
OnCalendar=*:0/10
Persistent=true

[Install]
WantedBy=timers.target
